# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx_error.log
stdout_logfile=/var/log/nginx_access.log
priority=100

[program:backend]
command=python scripts/auth_api_server.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/backend_error.log
stdout_logfile=/var/log/backend_access.log
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"
priority=200
startretries=3
startsecs=10

[program:frontend]
command=npm start
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/frontend_error.log
stdout_logfile=/var/log/frontend_access.log
environment=NODE_ENV="production",PORT="3000"
priority=300
startretries=3
startsecs=10 